# ws-client/systemd/market-data-listener.service.j2
[Unit]
Description=Binance market data WebSocket listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}/src
Environment=SYMBOL={{ symbol }}
Environment=STREAM_TYPE={{ stream_type }}
Environment=LOG_DIR={{ app_dir }}/logs
ExecStart={{ app_dir }}/venv/bin/python3 websocket_market_data_listener.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target